# deploy/gunicorn-staging.service.example
# Copy to /etc/systemd/system/gunicorn-staging.service (adjust paths and user).
# Then: sudo systemctl daemon-reload && sudo systemctl enable --now gunicorn-staging
#
# Match nginx upstream in deploy/nginx-staging.conf.template (127.0.0.1:8000).

[Unit]
Description=PGP Django API (staging) — gunicorn
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/pgp-main-web/pgp-main-web/backend
EnvironmentFile=/root/pgp-main-web/pgp-main-web/backend/.env
ExecStart=/root/pgp-main-web/pgp-main-web/backend/.venv/bin/gunicorn \
  --bind 127.0.0.1:8000 \
  --workers 3 \
  --timeout 120 \
  config.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
